# stow/i3/.config/systemd/user/gammastep.service.d/display.conf
[Service]
Environment=DISPLAY=:0
# Gammastep is a Wayland tool
Environment=WAYLAND_DISPLAY=wayland-1
