# stow/i3/.config/systemd/user/lxpolkit.service.d/display.conf
[Service]
Environment=DISPLAY=:1
Environment=XAUTHORITY=/run/user/1000/gdm/Xauthority
